# deploy/bookstore-scraper.service
[Unit]
Description=Bookstore Scraper API Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/bookstore-scraper
ExecStart=/opt/bookstore-scraper/.venv/bin/python -m src.cli serve
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
